# deploy/launchd/com.macmin.dashboard.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.macmin.dashboard.worker</string>
	<key>WorkingDirectory</key>
	<string>{{REPO_ROOT}}</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>DASHBOARD_CONFIG_PATH</key>
		<string>{{REPO_ROOT}}/config/config.yaml</string>
		<key>DASHBOARD_DB_PATH</key>
		<string>{{REPO_ROOT}}/data/fleet.db</string>
		<key>DASHBOARD_AUDIT_INTERVAL_SEC</key>
		<string>300</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>{{WORKER_BIN}}</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>{{LOG_DIR}}/worker.log</string>
	<key>StandardErrorPath</key>
	<string>{{LOG_DIR}}/worker.err</string>
</dict>
</plist>
